# deploy/hdaq/odin_control.conf
[program:odin_control]
directory=/hxt_sw/src/hexitec-detector/control/
command=bash -c "source /hxt_sw/install/v-abort/bin/activate && odin_control --config config/hexidaq_control.cfg --http_addr 192.168.0.27" -DFOREGROUND

user=hexiuser
environment=HOME="/home/hexiuser",USER="hexiuser"
autostart=true
autorestart=true
stdout_logfile=/tmp/odin_control.out
