# deploy/hdaq/legacy/odin_data_fr_1_v100.conf
[program:odin_data_fr_1_v100]
directory=/aeg_sw/work/users/ckd27546/develop/projects/hexitec/hxt_sw/src/install/
command=bash -c "sleep 4 && bin/frameReceiver -d 0 --json_file config/data/fr_hexidaq_config_0.json --ctrl tcp://0.0.0.0:5014 --logconfig config/data/fr_log4cxx.xml --rxbuffer 500000000" -DFOREGROUND
user=hexiuser
environment=HOME="/home/hexiuser",USER="hexiuser"
autostart=true
autorestart=true
stderr_logfile=/tmp/odin_data_fr_1_v100.err
stdout_logfile=/tmp/odin_data_fr_1_v100.out
