# deploy/hdaq/legacy/odin_control_v100.conf
[program:odin_control_v100]
directory=/hxt_sw/src/hexitec-detector/control
command=bash -c "source /hxt_sw/src/install/v-hdaq-1.0.0-ed/bin/activate && odin_control --config config/hexitec_dev_pc.cfg --logging=debug --debug_mode=1" -DFOREGROUND
user=hexiuser
environment=HOME="/home/hexiuser",USER="hexiuser"
autostart=true
autorestart=true
stderr_logfile=/tmp/odin_control_v100.err
stdout_logfile=/tmp/odin_control_v100.out
