# deploy/hdaq/legacy/odin_data_fp_2_v100.conf
[program:odin_data_fp_2_v100]
directory=/aeg_sw/work/users/ckd27546/develop/projects/hexitec/hxt_sw/src/install/
command=bash -c "bin/frameProcessor -d 0 --json_file config/data/fp_hexitec_config_1.json --ctrl tcp://127.0.0.1:5104 --logconfig config/data/fp_log4cxx.xml" -DFOREGROUND
user=hexiuser
environment=HOME="/home/hexiuser",USER="hexiuser"
autostart=true
autorestart=true
stderr_logfile=/tmp/odin_data_fp_2_v100.err
stdout_logfile=/tmp/odin_data_fp_2_v100.out
